# etc/systemd/system/radio.service
[Unit]
Description=Internet Radio

[Service]
ExecStart=/home/pi/.local/bin/radio

[Install]
WantedBy=multi-user.target